# .config/systemd/user/reloadHypr.service
[Unit]
Description=Reload hyprland

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprctl reload

[Install]
WantedBy=default.target
